# packages/web/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
priority=5

[program:inertia-ssr]
command=php artisan inertia:start-ssr
autostart=true
autorestart=true
priority=15